# SBFspotUploadDaemon/SBFspotUpload.service.default
[Unit]
Description=SBFspot upload daemon
After=network.target

[Service]
User=root
Type=forking
TimeoutStopSec=10
ExecStart=/usr/local/bin/sbfspot.3/SBFspotUploadDaemon  2>&1> /dev/null

Restart=on-failure
RestartSec=42s
SuccessExitStatus=SIGKILL

[Install]
WantedBy=default.target
